# deploy/services/cloudflare-ddns.timer.j2
[Unit]
Description=Timer for Cloudflare DDNS
Requires=cloudflare-ddns.service

[Timer]
Unit=cloudflare-ddns.service
OnUnitInactiveSec=300

[Install]
WantedBy=timers.target
